# streamer.conf
description "Twitter stream service"
author "Mark"

start on (local-filesystems and net-device-up)
stop on runlevel [06]

respawn
exec /home/{username}/streamapp/app.py
